# playbooks/files/certbot.service
[Unit]
Description=Let's Encrypt Certificate Renewal

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew --agree-tos
